# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

; Bootstrap: cache config, run migrations, seed categories, then start the web server
[program:web]
command=/bin/sh -c "php artisan config:cache && php artisan route:cache && php artisan view:cache && php artisan storage:link && php artisan migrate --force && php artisan db:seed --class=CategorySeeder --force && php artisan db:seed --class=SubscriptionPlanSeeder --force && php artisan serve --host=0.0.0.0 --port=8000"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; WebSocket server — required for real-time notifications
[program:reverb]
command=php artisan reverb:start --host=0.0.0.0 --port=8090
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; Queue worker — processes background jobs and notification emails
[program:queue]
command=php artisan queue:work --sleep=3 --tries=3 --max-time=3600
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
